# private_dot_config/supercronic/crontab
# Durability saves for the work tmux server (supercronic reads this file).
# Staggered by a minute: run concurrently, a tmux-snapshot's queries race the
# resurrect save running under tmux-save and the save silently produces no file.
*/5 * * * * $HOME/.local/bin/tmux-save
1-59/5 * * * * $HOME/.local/bin/tmux-snapshot

# Time tracking: auto-fill Toggl gaps hourly (no-op when no gaps detected).
0 * * * * $HOME/.local/bin/fill-toggl-cron.sh

# Prune opencode sessions older than 8 weeks, daily.
30 4 * * * $HOME/.local/bin/opencode-prune

# Bound logs/saves/data with no built-in rotation, daily (resurrect saves keep
# 5000, append logs 1 GiB ceiling, tt events keep ~2 months).
0 4 * * * $HOME/.local/bin/rotate-logs
